# docker/supervisor/scanner.ini
[program:scanner]
command = bin/throttle %(ENV_SCAN_INTERVAL)s bin/console webf:flysystem-failover:sync
autorestart = true
directory = /opt/app
user = www-data

stdout_events_enabled = true
stderr_events_enabled = true
