# docker/supervisor/scanner.ini
[program:scanner]
command = bin/scanner
autorestart = true
directory = /opt/app

stdout_events_enabled = true
stderr_events_enabled = true
