# docker/supervisor/workers.ini
[program:worker]
command = bin/console webf:flysystem-failover:process-message --limit=1
process_name = %(program_name)s_%(process_num)s
directory = /opt/app
user = www-data

numprocs = %(ENV_NB_WORKERS)s
autorestart = true

stdout_events_enabled = true
stderr_events_enabled = true
